# archlive/airootfs/etc/systemd/system/yantra-daemon.service
[Unit]
Description=YantraOS Kriya Loop Daemon
After=network.target

[Service]
Type=simple
User=yantra_user
Group=yantra_user
RuntimeDirectory=yantra
RuntimeDirectoryMode=0755
ExecStart=/opt/yantra/venv/bin/python3 /opt/yantra/core/daemon.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
