# deploy/systemd/yantra-provision-secrets.service
[Unit]
Description=YantraOS First-Boot Secret Provisioning
Before=yantra.service yantra-telegram.service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=infinity
StartLimitBurst=90

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
ExecStart=/opt/yantra/venv/bin/python3 /opt/yantra/scripts/provision_runtime_secrets.py
RemainAfterExit=yes
Restart=on-failure
RestartSec=10
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/yantra
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_FOWNER CAP_CHOWN
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
