# archlive/airootfs/usr/lib/systemd/system/emergency.service.d/override.conf
[Service]
# SECURITY: SYSTEMD_SULOGIN_FORCE=1 removed — emergency shells must
# demand root password authentication. No physical bypass permitted.
