# archlive/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
# ────────────────────────────────────────────────────────────────────────
# YantraOS — Kiosk TUI Auto-Login Override
# Target: /etc/systemd/system/getty@tty1.service.d/override.conf
#
# Automatically logs in yantra_user on TTY1 and launches tui_shell.py
# as the primary interface. The underlying Linux terminal is completely
# masked — the user never sees a root shell.
#
# The -o flag passes login options:
#   -p  preserve environment
#   -f  skip authentication (auto-login)
#
# ExecStartPost launches the TUI via the Python virtualenv.
# ────────────────────────────────────────────────────────────────────────

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin yantra_user --noclear %I $TERM
Type=idle

[Service]
ExecStartPost=-/bin/sh -c 'sleep 1 && su - yantra_user -c "/opt/yantra/venv/bin/python /opt/yantra/core/tui_shell.py" &'
