# archlive/airootfs/etc/systemd/system/yantra.service.d/env.conf
[Service]
EnvironmentFile=/etc/yantra/.env.secrets
